# defaults/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisord]
nodaemon=true
logfile=/dev/fd/1
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root
childlogdir=/tmp
; Make sure children inherit a sane umask
umask=022

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ============================
; Global defaults (manually repeated per program because supervisord
; does not support true inheritance for these fields)
; - Route logs to Docker
; - No log rotation
; ============================

; ============================
; Xvfb
; ============================
[program:xvfb]
command=/bin/sh -lc 'W=${SCREEN_WIDTH:-1920}; H=${SCREEN_HEIGHT:-1080}; DPI=${SCREEN_DPI:-96}; EXTRA=${XVFB_EXTRA_ARGS:-}; exec /usr/bin/Xvfb "${DISPLAY}" -screen 0 ${W}x${H}x24 -dpi ${DPI} +extension RANDR +extension DPMS -nolisten tcp -ac -noreset ${EXTRA}'
autostart=true
autorestart=true
priority=10
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s",SCREEN_WIDTH="%(ENV_SCREEN_WIDTH)s",SCREEN_HEIGHT="%(ENV_SCREEN_HEIGHT)s",SCREEN_DPI="%(ENV_SCREEN_DPI)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; dbus system bus (fix xfdesktop "Failed to get system bus")
; ============================
[program:dbus_system]
command=/bin/sh -lc 'mkdir -p /run/dbus && chmod 755 /run/dbus; [ -s /etc/machine-id ] || dbus-uuidgen --ensure=/etc/machine-id; exec /usr/bin/dbus-daemon --system --nofork --nopidfile'
autostart=true
autorestart=true
priority=15
user=root
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; wineserver (keepalive)
; ============================
[program:wineserver]
command=/bin/sh -lc 'exec /data/conf/scripts/run_wineserver.sh'
autostart=true
autorestart=true
priority=18
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
startsecs=1

; ============================
; XFCE session
; ============================
[program:xfce]
command=/bin/sh -lc 'export HOME="%(ENV_HOME)s"; D="${DISPLAY#:}"; for i in $(seq 1 400); do [ -S "/tmp/.X11-unix/X${D}" ] && break; sleep 0.05; done; exec /data/conf/scripts/run_xfce.sh'
autostart=true
autorestart=true
priority=20
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s",WALLPAPER="%(ENV_WALLPAPER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
startsecs=3

; ============================
; SSHD
; ============================
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
priority=22
user=root
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; autocutsel (clipboard sync)
; ============================
[program:autocutsel_clipboard]
command=/bin/sh -lc 'D="${DISPLAY#:}"; for i in $(seq 1 200); do [ -S "/tmp/.X11-unix/X${D}" ] && break; sleep 0.05; done; exec /usr/bin/autocutsel -selection CLIPBOARD'
autostart=true
autorestart=true
startsecs=1
priority=25
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; x11vnc (waits for correct X socket)
; ============================
[program:x11vnc]
command=/bin/sh -lc 'D="${DISPLAY#:}"; for i in $(seq 1 200); do [ -S "/tmp/.X11-unix/X${D}" ] && break; sleep 0.05; done; exec /usr/bin/x11vnc -display "${DISPLAY}" -rfbport 5900 -forever -noprimary -shared -afteraccept /data/conf/scripts/onconnect.sh -gone /data/conf/scripts/ondisconnect.sh -nopw -xkb -noxdamage -nonap -wait 10 -deferupdate 10 -repeat -skip_lockkeys'
autostart=true
autorestart=true
priority=30
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; websockify
; ============================
[program:websockify]
command=/usr/bin/websockify 127.0.0.1:6080 127.0.0.1:5900
autostart=true
autorestart=true
priority=40
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; toolbar_api
; ============================
[program:toolbar_api]
command=/usr/bin/python3 /data/conf/scripts/toolbar_api.py
autostart=true
autorestart=true
priority=45
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s",TOOLBAR_API_PORT="9001",WM_SERVICE="xfce",DESKTOP_SERVICE="none",KIOSK_SCRIPT="/data/conf/scripts/kiosk_mode.sh",TOOLBAR_API_STATIC_DIR="/data/conf/www",SUPERVISOR_STATUS_EMPTY_OK="0"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; nginx
; ============================
[program:nginx]
command=/bin/sh -lc 'mkdir -p /tmp/nginx/client_body /tmp/nginx/proxy /tmp/nginx/fastcgi /tmp/nginx/uwsgi /tmp/nginx/scgi; /usr/sbin/nginx -t -c /data/conf/nginx/nginx.conf -g "pid /tmp/nginx.pid;" && exec /usr/sbin/nginx -c /data/conf/nginx/nginx.conf -g "daemon off; pid /tmp/nginx.pid;"'
autostart=true
autorestart=true
priority=50
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

; ============================
; startup_script_once
; ============================
[program:startup_script_once]
command=/data/conf/scripts/prefix-log.sh /data/conf/scripts/run_startup_script_once.sh  
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
priority=999
user=%(ENV_USER)s
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s",STARTUP_SCRIPT="%(ENV_STARTUP_SCRIPT)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

